# vlc/source/asselgrabber.conf
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

respawn

exec su - assel -c '/home/assel/asseln/asselgrabber.sh 2>/var/log/asseln/grabber.log'
